# scripts/lr-autotag.service
[Unit]
Description=LR-AutoTag Backend
After=network.target postgresql.service

[Service]
Type=exec
User=lr-autotag
WorkingDirectory=/opt/lr-autotag
EnvironmentFile=/opt/lr-autotag/.env
ExecStart=/opt/lr-autotag/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
